# sriov-createvf-1/etc/systemd/system/sriov-createvf.service
[Unit]
Description=Create Virtual Function devices
After=network-online.target

[Service]
EnvironmentFile=-/etc/sysconfig/sriov-createvf
Type=oneshot
RemainAfterExit=yes
ExecStart=-/usr/bin/sriov-createvf ${DEVICES}

[Install]
WantedBy=multi-user.target
